# .config/systemd/user/udiskie.service
[Unit]
Description=Udiskie USB Drives Automount

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target

